# device-httpd-el9/src/main/resources/etc/profile.d/device-status-httpd.sh
#!/bin/bash

if /usr/bin/tty -s
then
  /usr/libexec/device-status/httpd
fi

